# deploy/docker/systemd/cloudport-docker-cleanup.service
[Unit]
Description=CloudPort - limpeza de imagens e cache Docker
Documentation=https://github.com/diogo2806/CloudPort
Requires=docker.service
After=docker.service
ConditionPathExists=/var/run/docker.sock

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/cloudport-docker-cleanup
ExecStart=/usr/local/sbin/cloudport-docker-cleanup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=30min

[Install]
WantedBy=multi-user.target